# deploy/delta-bot.service
[Unit]
Description=Delta Exchange Trading Bot
After=network.target

[Service]
Type=simple
User=YOUR_VM_USER
WorkingDirectory=YOUR_WORKING_DIR
EnvironmentFile=YOUR_WORKING_DIR/.env
# This will be replaced by setup.sh with the actual venv python path
ExecStart=YOUR_PYTHON_EXEC main_v2.py
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=delta-bot

[Install]
WantedBy=multi-user.target
